# debian/loolwsd.cron.d
0 0 */1 * * root find /var/cache/loolwsd -type f -a -atime +10 -delete
